# supervisor/horizon.conf
[program:laravel-horizon]
command=/usr/local/bin/php /var/www/html/artisan horizon
autostart=true
autorestart=true
